# autostart/.config/autostart/fusuma.desktop
[Desktop Entry]
Name=Fusuma
GenericName=fusuma
Exec=/usr/local/bin/fusuma
Type=Application
X-GNOME-Autostart-Delay=30
